# extra/cron.monthly
#!/bin/sh

/opt/mojo-acme-server/mojo-acme-server acme automate -T nginx_default -l http://*:8928 -o proxy_pass=http://127.0.0.1:3000 example.com www.example.com
/opt/mojo-acme-server/mojo-acme-server acme automate -T nginx_default -l http://*:8928 -o proxy_pass=http://127.0.0.1:3000 example1.com www.example1.com
/etc/init.d/nginx testconfig && systemctl restart nginx.service
